# debian/upstart/nfcd.conf
description "NFC daemon"

start on started dbus and android
stop on stopping dbus

exec /usr/sbin/nfcd -p /usr/lib/nfcd/plugins
